# roles/smartctl_exporter/templates/smartctl_exporter.service.j2
[Unit]
Description=SmartCTL Exporter
After=network.target

[Service]
User=root
ExecStart=/usr/local/bin/smartctl_exporter-0.14.0.linux-amd64/smartctl_exporter

[Install]
WantedBy=multi-user.target

